# reach/external/mpu_broadcast/mpu_broadcast.service
[Unit]
Description="MPU raw data broadcaster."
After=reach-setup.service usb0.network

[Service]
ExecStart=
ExecStart=/usr/bin/mpu_broadcast
Type=simple
